# home/Library/LaunchAgents/renders.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>local.renders</string>
	<key>LowPriorityIO</key>
	<true/>
	<key>ProcessType</key>
	<string>Background</string>
	<key>ProgramArguments</key>
	<array>
		<string>zsh</string>
		<string>-ci</string>
		<string>$HOME/.config/scripts/renders.sh</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/tmp/local.renders.stderr</string>
	<key>StandardOutPath</key>
	<string>/tmp/local.renders.stdout</string>
	<key>WatchPaths</key>
	<array>
		<string>~/Library/Mobile Documents/com~apple~CloudDocs/Ableton/Exports</string>
	</array>
</dict>
</plist>
